# system/hbank-payment-plans.timer
[Unit]
Description=Run H-Bank Payment Plans Daily

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
